# apps/operation/community-monitor/community-monitor.service
[Unit]
Description=Community Model Monitor Loop
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/monitor
EnvironmentFile=/home/ubuntu/monitor/.env
ExecStart=/home/ubuntu/monitor/loop.sh
Restart=always
RestartSec=900
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
