# social/reddit/systemd/reddit_bot.timer
[Unit]
Description=Run Reddit Bot Deploy Daily

[Timer]
OnCalendar=*-*-* 10:30:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target