# image.pollinations.ai/nunchaku/flux-server.service
[Unit]
Description=FLUX Schnell Nunchaku Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/pollinations/image.pollinations.ai/nunchaku
Environment=PATH=/home/ubuntu/pollinations/image.pollinations.ai/nunchaku/nunchaku_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/home/ubuntu/.env
ExecStart=/home/ubuntu/pollinations/image.pollinations.ai/nunchaku/nunchaku_env/bin/python server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
