# operations/systemd/text-pollinations.service
[Unit]
Description=Pollinations Text API Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/pollinations/text.pollinations.ai
Environment="NODE_ENV=production"
Environment="PORT=16385"
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=text-pollinations

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
